# systemd/.config/systemd/user/home-backup.timer
[Unit]
Description=Timer to backup user home dir

[Timer]
OnCalendar=*-*-* 17:00

[Install]
WantedBy=timers.target
